# Kubernetes/server/srcs/influxdb_old/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:influxdb]
command = sh -c "influxd  && kill -s SIGTERM $(cat supervisord.pid)"
